# ubuntu-plc-sender.service
[Unit]
Description=Ubuntu PLC API Sender
After=network.target
Wants=network.target

[Service]
Type=simple
User=plcuser
Group=plcuser
WorkingDirectory=/opt/ubuntu-plc-sender
ExecStart=/opt/ubuntu-plc-sender/UbuntuPlcApiSender
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ubuntu-plc-sender

# Environment variables
Environment=DOTNET_ENVIRONMENT=Production
Environment=ASPNETCORE_URLS=

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/ubuntu-plc-sender

# Resource limits
LimitNOFILE=65536
MemoryMax=512M

# Restart policy
StartLimitIntervalSec=300
StartLimitBurst=5

[Install]
WantedBy=multi-user.target
